# G53IDS/.config/gunicorn.conf
[program:gunicorn]
command=/root/.virtualenvs/env-g53ids/bin/gunicorn run:app -b localhost:5000
directory=/home/www/G53IDS
